# scripts/install/launchd/com.forgewire.runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.forgewire_fabric.runner</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/python3</string>
    <string>-m</string>
    <string>forgewire_fabric.cli</string>
    <string>runner</string>
    <string>start</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>FORGEWIRE_HUB_URL</key>
    <string>http://hub.local:8765</string>
    <key>FORGEWIRE_HUB_TOKEN_FILE</key>
    <string>/usr/local/etc/forgewire/hub.token</string>
    <key>FORGEWIRE_RUNNER_WORKSPACE_ROOT</key>
    <string>/usr/local/var/forgewire/workspace</string>
    <key>FORGEWIRE_RUNNER_TAGS</key>
    <string>macos,python:3.11</string>
    <key>FORGEWIRE_RUNNER_SCOPE_PREFIXES</key>
    <string>src/,tests/</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/usr/local/var/forgewire</string>

  <key>UserName</key>
  <string>_forgewire</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>NetworkState</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/var/log/forgewire/runner.out.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/forgewire/runner.err.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
