# scripts/install/launchd/com.forgewire.hub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ForgeWire Hub launchd daemon.

  Install:
    sudo cp com.forgewire_fabric.hub.plist /Library/LaunchDaemons/
    sudo chown root:wheel /Library/LaunchDaemons/com.forgewire_fabric.hub.plist
    sudo chmod 0644       /Library/LaunchDaemons/com.forgewire_fabric.hub.plist
    sudo launchctl load -w /Library/LaunchDaemons/com.forgewire_fabric.hub.plist

  Edit the ProgramArguments path to point at your Python interpreter and
  adjust the EnvironmentVariables block (token file, port, db path).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.forgewire_fabric.hub</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/python3</string>
    <string>-m</string>
    <string>forgewire_fabric.cli</string>
    <string>hub</string>
    <string>start</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>8765</string>
    <string>--db-path</string>
    <string>/usr/local/var/forgewire/hub.sqlite3</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>FORGEWIRE_HUB_TOKEN_FILE</key>
    <string>/usr/local/etc/forgewire/hub.token</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/usr/local/var/forgewire</string>

  <key>UserName</key>
  <string>_forgewire</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>NetworkState</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/var/log/forgewire/hub.out.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/forgewire/hub.err.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
