# python/forgewire_fabric/_installer_assets/systemd/forgewire-runner.service
[Unit]
Description=ForgeWire Runner (claims signed remote dispatches)
Documentation=https://github.com/DigitalHallucinations/forgewire-fabric
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# --- Edit these for your install ---------------------------------------------
User=forgewire
Group=forgewire
WorkingDirectory=/var/lib/forgewire

Environment=PYTHON=/usr/bin/python3
Environment=FORGEWIRE_HUB_URL=http://hub.local:8765
Environment=FORGEWIRE_HUB_TOKEN_FILE=/etc/forgewire/hub.token
Environment=FORGEWIRE_RUNNER_WORKSPACE_ROOT=/var/lib/forgewire/workspace
Environment=FORGEWIRE_RUNNER_TAGS=linux,python:3.11
Environment=FORGEWIRE_RUNNER_SCOPE_PREFIXES=src/,tests/
Environment=FORGEWIRE_RUNNER_MAX_CONCURRENT=1
Environment=FORGEWIRE_RUNNER_POLL_INTERVAL=5
# -----------------------------------------------------------------------------

ExecStart=/bin/sh -c '${PYTHON} -m forgewire_fabric.cli runner start \
    --workspace-root ${FORGEWIRE_RUNNER_WORKSPACE_ROOT} \
    --tags ${FORGEWIRE_RUNNER_TAGS} \
    --scope-prefixes ${FORGEWIRE_RUNNER_SCOPE_PREFIXES}'

Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/forgewire /var/log/forgewire
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged

StandardOutput=journal
StandardError=journal
SyslogIdentifier=forgewire-runner

[Install]
WantedBy=multi-user.target
